# cv_pipeline/meza-occupancy-detector@.service
[Unit]
Description=MEZA occupancy CV pipeline (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/meza/cv_pipeline
Environment=CAMERA_ID=%i
Environment=SUPABASE_URL=https://your-project.supabase.co
Environment=SUPABASE_SERVICE_KEY=<service role key>
ExecStart=/usr/bin/python3 /opt/meza/cv_pipeline/occupancy_detector.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target

# This is a systemd template unit (note the %i / @ in the filename) - one
# instance per camera on this device:
#   sudo cp meza-occupancy-detector.service /etc/systemd/system/meza-occupancy-detector@.service
#   sudo systemctl enable --now meza-occupancy-detector@<camera-id>
# Replace the SUPABASE_URL/SUPABASE_SERVICE_KEY placeholders above (or
# better, use an EnvironmentFile= line pointing at a root-only-readable
# file instead of inlining the service key in this unit).
